# files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI
uci set pbuf.opt.memory_profile=auto
uci set firewall.@defaults[0].flow_offloading=0
uci set network.globals.packet_steering=0
uci set ecm.@general[0].enable_bridge_filtering=1
EOI
